# systemd/tinyagentos-host-firewall.path
[Unit]
Description=React to Docker appearing after taOS install; ensure DOCKER-USER contains our accept rules
After=tinyagentos-host-firewall.service

[Path]
# Fires whenever /var/run/docker.pid appears — the cleanest signal that
# Docker has started or restarted.  The associated service is idempotent
# so multiple triggers are harmless.
PathExists=/var/run/docker.pid
Unit=tinyagentos-host-firewall.service

[Install]
WantedBy=multi-user.target
