# systemd/tinyagentos-host-firewall.service
[Unit]
Description=tinyagentos host firewall rules (allow incus bridge through docker DROP)
After=network-online.target docker.service incus.service
Before=tinyagentos.service
Wants=network-online.target
# Skip on pure-nftables systems that don't ship the iptables legacy wrapper.
# Either location is acceptable; the glob matches either.
ConditionPathExistsGlob=/usr/sbin/iptables /usr/bin/iptables

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/tinyagentos/scripts/host-firewall-up.sh
ExecStop=/opt/tinyagentos/scripts/host-firewall-down.sh

[Install]
WantedBy=multi-user.target
