# systemd/tinyagentos-host-firewall.service
[Unit]
Description=tinyagentos host firewall rules (allow incus bridge through docker DROP)
After=network-online.target docker.service incus.service
Before=tinyagentos.service
Wants=network-online.target
# Skip on pure-nftables systems that don't ship the iptables legacy wrapper.
# Multiple ConditionPathExists directives are AND-combined; the "|" prefix
# on each turns them into OR so the unit runs if iptables exists at EITHER
# location. A single-line space-separated glob was silently evaluating to
# false and skipping the unit on every boot.
ConditionPathExists=|/usr/sbin/iptables
ConditionPathExists=|/usr/bin/iptables

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/tinyagentos/scripts/host-firewall-up.sh
ExecStop=/opt/tinyagentos/scripts/host-firewall-down.sh

[Install]
WantedBy=multi-user.target
